# archivator.cron
15 4 * * * root /opt/replicator/archivator >> /var/log/archivator.log
