# replicator.cron
*/5 * * * * root /opt/replicator/replicator.rsync

